# molecule/website/roles/website/templates/crontab
*/5 * * * * flock /tmp/update-website /srv/update-website.sh >> /var/log/update-website.log 2>&1
