# molecule/packages/roles/packages/templates/crontab
*/5 * * * * flock /tmp/update-packages /srv/update-packages.sh >> /var/log/update-packages.log 2>&1
