# molecule/packages/roles/packages/templates/crontab
*/5 * * * * flock /tmp/update-packages /srv/update-packages.sh
